# kyoko.container
[Container]
ContainerName=kyoko
Image=kyoko-client
PublishPort=8881:80

[Service]
Restart=always

[Install]
WantedBy=default.target
